# platforms/zcu102_arm/images/aarch64-xilinx-linux/usr/lib/systemd/user/at-spi-dbus-bus.service
[Unit]
Description=Accessibility services bus

[Service]
Type=dbus
BusName=org.a11y.Bus
ExecStart=/usr/libexec/at-spi-bus-launcher
